# scripts/systemd/advisor-ledger.timer
[Unit]
Description=Run Advisor Ledger pipeline every 10 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=15s
Persistent=true
Unit=advisor-ledger.service

[Install]
WantedBy=timers.target
